# jira-sync.service
[Unit]
Description=Jira Sync Application
After=network.target postgresql.service

[Service]
User=hubert
WorkingDirectory=/home/hubert/dataappV3
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/hubert/dataappV3/app.py
Restart=always
RestartSec=10

# Ensure proper file permissions
UMask=0002

# Give the service time to save config and clean up
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
